# 0Linux/initialisation-systeme/rc.local
#!/bin/sh
#
# /etc/rc.d/rc.local : Initialisation locale du système.
#
# Placez ici vos commandes personnalisées à lancer à la fin de l'amorçage.
# Pour exécuter des commandes à l'arrêt ou au redémarrage de la machine,
# vous pouvez créer un fichier /etc/rc.d/rc.local_shutdown et y placer
# les commandes désirées (mais n'oubliez pas de le rendre exécutable).
